# board/raspberrypi4-64/rootfs_overlay/etc/profile.d/shell.sh
#!/bin/sh

#------------------------------------------------------------------------

export LANG="ru_RU.UTF-8";

lRed()   { echo -n "\e[01;31m\]$@\e[00m\]"; }
lGreen() { echo -n "\e[01;32m\]$@\e[00m\]"; }
lBlue()  { echo -n "\e[01;34m\]$@\e[00m\]"; }

PS_TITLE="\h";
PS_PATH="$(lBlue "\w \$") ";

if [ "$PS1" ];
then
    case $(id -u) in
        0) export PS1="$(lRed      "\h") $(lBlue "\w # ")"; ;;
        *) export PS1="$(lGreen "\u@\h") $(lBlue "\w \$ ")"; ;;
    esac;
fi;

#------------------------------------------------------------------------

#------------------------------------------------------------------------
